# mmb-autoscroll.service
[Unit]
Description=mmb-autoscroll — Windows-style middle-button autoscroll (evdev grab -> uinput)
Documentation=https://wayland.freedesktop.org/libinput/doc/latest/scrolling.html
After=systemd-udevd.service
Wants=systemd-udevd.service
# Rate-limit directives belong in [Unit] (systemd ignores them in [Service],
# which is why the old unit logged an "Unknown key name" warning every start).
# With the wait-for-mouse fix the process no longer exit-loops, but keep the
# limit off so a genuine crash can't ever trip a start-limit lockout.
StartLimitIntervalSec=0

[Service]
Type=simple
# Tunables live here (optional file; "-" = ignore if missing). Edit it and
# `systemctl restart mmb-autoscroll` to retune without touching the unit.
EnvironmentFile=-/etc/mmb-autoscroll.conf
ExecStart=/usr/local/bin/mmb-autoscroll
Restart=always
RestartSec=2

# --- hardening -----------------------------------------------------------
# Needs only the input subsystem + the uinput node; deny everything else.
DevicePolicy=closed
DeviceAllow=/dev/uinput rw
DeviceAllow=char-input rw
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_UNIX
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
